# assets/selenium/mdc_debug_port.conf
[program:mdc_debug]
priority=100
command=/mdc_debug_port.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/dev/stdout
